# etc/profile.d/meeseeks.sh
#!/bin/bash

export MEESEEKS_HOME=/home/meeseeks

export PATH=$PATH:$MEESEEKS_HOME/bin
export PYTHONPATH=$MEESEEKS_HOME/lib

export MEESEEKS_CONF=''
